# nice/.config/systemd/user/mailpit.service
[Unit]
Description=Mailpit server

[Service]
ExecStart=/usr/bin/mailpit -d /var/lib/mailpit/mailpit.db
Restart=always
# Restart service after 10 seconds if service crashes
RestartSec=10
SyslogIdentifier=mailpit
User=mailpit
Group=mailpit

[Install]
WantedBy=multi-user.target
